# playbooks/templates/phd2.service.j2
[Unit]
Description=PHD2
Requires=Xvfb_phd2.service x11vnc_phd2.service
After=Xvfb_phd2.service x11vnc_phd2.service

[Service]
Environment="DISPLAY=:1"
ExecStart=/usr/local/bin/phd2

[Install]
