# playbooks/templates/x11vnc_phd2.service.j2
[Unit]
Description=VNC client for PHD2
Requires=Xvfb_phd2.service
After=Xvfb_phd2.service
StopWhenUnneeded=true

[Service]
ExecStart=/usr/bin/x11vnc -display :1

[Install]
