# playbooks/templates/Xvfb_phd2.service.j2
[Unit]
Description=Virtual X frame buffer for PHD2
StopWhenUnneeded=true

[Service]
ExecStart=/usr/bin/Xvfb :1 -screen 0 1200x800x16

[Install]
